# devilspie2/.config/autostart/devilspie2.desktop
[Desktop Entry]
Name=devilspie2
Exec=/usr/bin/devilspie2
Terminal=false
Type=Application
StartupNotify=false
